# supervisord/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autorestart=true
stdout_logfile=/var/log/redis/stdout.log
stderr_logfile=/var/log/redis/stderr.log
