# supervisord/conf.d/server.conf
[program:server]
command=python app.py
autostart=true
autorestart=true
directory=/app
user=nobody
stdout_logfile=/var/log/server/stdout.log
stderr_logfile=/var/log/server/stderr.log
